# launchd/com.shichangliao.ainews.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.shichangliao.ainews.sync</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/node</string>
      <string>/Users/shichangliao/Desktop/codes/MyAppStore/AINews/scripts/auto-sync-and-publish.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/shichangliao/Desktop/codes/MyAppStore/AINews</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/Users/shichangliao/Desktop/codes/MyAppStore/AINews/logs/sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/shichangliao/Desktop/codes/MyAppStore/AINews/logs/sync.error.log</string>
  </dict>
</plist>
