# launchd/com.ainews.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.ainews.sync</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/env</string>
      <string>node</string>
      <string>/Users/sc-claw/Desktop/Github/ai-news/scripts/auto-sync-and-publish.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/sc-claw/Desktop/Github/ai-news</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/Users/sc-claw</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>GIT_TERMINAL_PROMPT</key>
      <string>0</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/tmp/ai-news-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/ai-news-sync.error.log</string>
  </dict>
</plist>
